# deploy/cec.service
[Unit]
Description=Kiosk CEC

[Service]
User=pi
Group=pi
ExecStart=/home/pi/cec.sh

[Install]
WantedBy=default.target
